# Linux/etc/systemd/system/aws-mount-ephemeral.service.template
[Unit]
Description=Mount ephemeral drive to tmp

[Service]
ExecStart=/sbin/aws-mount-ephemeral /tmp
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
